# kernel/xpmem.service.in
[Unit]
Description=xpmem
After=basic.target
[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/sbin/modprobe --first-time xpmem
ExecStop=/sbin/rmmod xpmem
[Install]
WantedBy=multi-user.target
